# scripts/wii-controller.service
[Unit]
Description=Wii controller for robot

[Service]
Type=simple
ExecStart=/usr/bin/wii-controller-c

[Install]
WantedBy=multi-user.target
